# tmpfs-sync.sync.timer
[Unit]
Description=tmpfs-sync sync service

[Timer]
OnUnitActiveSec=2h
